# molecule/default/.molecule/roles/bngsudheer.ruby/templates/etc/profile.d/enableruby24.sh
#!/bin/bash
source /opt/rh/rh-ruby24/enable
export X_SCLS="`scl enable rh-ruby24 'echo $X_SCLS'`";
