# molecule/default/.molecule/roles/bngsudheer.ruby/templates/etc/profile.d/enableruby23.sh
#!/bin/bash
source /opt/rh/rh-ruby23/enable
export X_SCLS="`scl enable rh-ruby23 'echo $X_SCLS'`";
